# php82-moodle-unit/crontab
* * * * * /usr/local/bin/cron.sh >> /dev/stdout 2>&1
